# services/mirrors-index.path
[Unit]
Description=Mirrors index page update trigger

[Path]
PathModified=/srv/www/

[Install]
WantedBy=multi-user.target
