# services/mirrors-index.path
[Unit]
Description=Mirrors index page update trigger

[Path]
PathModified=/mnt/array0/pool/
PathModified=/mnt/array1/pool/
PathModified=/mnt/array1/exports/

[Install]
WantedBy=multi-user.target
