# services/mirrors-index.timer
[Unit]
Description=Update Mirrors Index Page Daily, If path file fails

[Timer]
OnCalendar=*-*-* 04:28:30

[Install]
WantedBy=timers.target
